# cacmin-bot.service
[Unit]
Description=CAC Admin Telegram Bot
Documentation=https://github.com/Cordtus/cacmin-bot
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=cacmin-bot
Group=cacmin-bot
WorkingDirectory=/opt/cacmin-bot

# Load environment variables from .env file
EnvironmentFile=/opt/cacmin-bot/.env

# Main process
ExecStart=/usr/bin/yarn start

# Restart configuration
Restart=always
RestartSec=10

# Environment
Environment=NODE_ENV=production

# Logging configuration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cacmin-bot

# Security hardening (optional, adjust as needed for LXC)
NoNewPrivileges=true
PrivateTmp=true

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
